# system/udev_rules/keyboard.rules
ACTION=="add", SUBSYSTEMS=="usb", ENV{DISPLAY}=":0", ENV{HOME}="/home/tyler", RUN+="/usr/bin/setxkbmap -v 10 -option caps:ctrl_modifier"
